# support/cron.d/IPStat
# run Watch script
*/1 * * * * root /opt/IPStat/sbin/IPStat-watch 1>/dev/null 2>/dev/null

